# systemd/.config/systemd/user/nc-drop.service
[Unit]
Description=Nextcloud drop folder uploader
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=%h/bin/nc-drop.sh
Restart=always
RestartSec=3

# Silence everything
StandardOutput=null
StandardError=null

[Install]
WantedBy=default.target

